# supervisord-tomcat.conf
[program:tomcat]
command=/start-tomcat.sh
numprocs=1
autostart=true
autorestart=false
user=root
